# etl/pg_to_es/crontab
0 0 * * * bash -l -c /app/run.sh >> /var/log/cron.log 2>&1
